# systemd/ogm_pi.service
# OGM_slave_pi daemon service (pair with ogm_pi.socket for activation).
[Unit]
Description=OGM Pi Modbus RTU + IPC Daemon
Requires=ogm_pi.socket
After=ogm_pi.socket
StartLimitIntervalSec=60
StartLimitBurst=5

[Service]
Type=simple
User=ogm_pi
Group=ogm
WorkingDirectory=/opt/OGM_slave_pi
Environment=PYTHONUNBUFFERED=1
ExecStart=/opt/OGM_slave_pi/.venv/bin/python -m ogm_pi.daemon --config /etc/ogm_pi/ogm_pi.yaml
Restart=always
RestartSec=2
TimeoutStopSec=5
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only
ProtectControlGroups=yes
ProtectKernelLogs=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
SupplementaryGroups=gpio dialout

[Install]
WantedBy=multi-user.target
